# mbus2mqtt/rootfs/etc/services.d/mbus2mqtt/finish
#!/bin/bash
bashio::log.info "M-Bus to MQTT stopped."
